# src/github.com/andrew-bodine/release/cmd/raspi-release-agent/raspi-release-agent@.service
[Unit]
Description=Raspi Release Agent
Documentation=https://github.com/andrew-bodine/raspi-workspace

[Service]
ExecStart=/opt/raspi/bin/raspi-release-agent

[Install]
WantedBy=multi-user.target
